# etc/systemd/system/avahi-daemon.service.d/override.conf
[Service]
ExecStartPost=-/usr/bin/wlanpi-avahi-service-update